# systemd/aria.target
# Copyright (c) 2025-present K. S. Ernest (iFire) Lee
# SPDX-License-Identifier: MIT

[Unit]
Description=Aria Character Core Complete Stack
Documentation=https://github.com/fire/aria-character-core
After=network.target
Wants=aria-cockroachdb.service aria-openbao.service aria-seaweedfs.service aria-app.service

[Install]
WantedBy=multi-user.target
